# raspicam/ansible/roles/cam/templates/etc/systemd/system/home-pi-cam.mount
[Mount]
What={{disk_path}}
Where=/home/pi/cam
Type=ext4
Options=defaults,noatime

[Install]
WantedBy=multi-user.target
